# .config/systemd/user/battery.service
[Unit]
Description=Battery low

[Service]
ExecStart=/home/nikita/.config/systemd/user/battery.sh

[Install]
WantedBy=default.target
