# systemd/gpu_undervolt.service
[Unit]
Description=GPU configuration

[Service]
ExecStart=gpu_undervolt.sh

[Install]
WantedBy=default.target
