# systemd/gpu_fan_curve.service
[Unit]
Description=GPU Fan Curve

[Service]
ExecStart=gpu_fan_curve.sh

[Install]
WantedBy=default.target
