# gentoo/user/init.d/polkit-gnome
#!/sbin/openrc-run

supervisor=supervise-daemon
command=/usr/libexec/polkit-gnome-authentication-agent-1
