# gentoo/user/init.d/shadercache
#!/sbin/openrc-run

start() {
  mkdir -p "/tmp/${USER}/shadercache"
}
