# gentoo/user/init.d/slstatus
#!/sbin/openrc-run

supervisor=supervise-daemon
command=/usr/bin/slstatus
