# gentoo/user/init.d/transmission
#!/sbin/openrc-run

supervisor=supervise-daemon
command=/usr/bin/transmission-daemon
command_args_foreground=--foreground
